# whisper-client.service
[Unit]
Description=Whisper Speech-to-Text Client
After=graphical-session.target

[Service]
Type=simple
Environment=PYTHONPATH=/home/mike/dev/local-stt-tool/venv/lib/python3.10/site-packages
ExecStart=/home/mike/dev/local-stt-tool/venv/bin/python /home/mike/dev/local-stt-tool/whisper_hotkey.py
Restart=always
RestartSec=5

[Install]
WantedBy=graphical-session.target